# app/resources/systemd/salt-update-grains.timer
[Unit]
Description=Update salt grains

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target